# deploy/sketch.service
[Unit]
Description=Sketch Discord bot (containerised)
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
EnvironmentFile=/etc/sketch/env
# Pull+extract can exceed systemd's 90s default on cold layers; give it room
# so a slow first-boot pull doesn't fail the unit. The inner `timeout 240`
# bounds the pull itself so a hung registry falls through to the cached image
# before systemd kills the whole start.
TimeoutStartSec=300
# Best-effort pull (leading "-"): if Artifact Registry is unreachable on boot,
# fall back to the locally cached image rather than failing the service.
ExecStartPre=-/usr/bin/timeout 240 /usr/bin/docker pull ${SKETCH_IMAGE}
# `--rm` cleans up the stopped container; systemd is the supervisor, not docker
# (so `--restart=no` would be the docker default anyway). Hardening flags here
# replace the host-level seccomp/AppArmor directives the previous unit used.
ExecStart=/usr/bin/docker run --rm --name sketch \
  --read-only --tmpfs /tmp \
  --cap-drop=ALL \
  --security-opt=no-new-privileges \
  --pids-limit=256 \
  -e DEV_GUILD_ID \
  -e GUILD_CONFIG_JSON \
  -e GCP_PROJECT \
  -e SKETCH_SECRET_DISCORD_TOKEN \
  ${SKETCH_IMAGE}
ExecStop=/usr/bin/docker stop sketch
Restart=always
RestartSec=10
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
